# docker/crontab
# Scheduled tasks executed by the seo_cron container (cron -f).
#
# telegram_sender: ticks every minute but is gated by the
# seo_settings.scheduled_publish_enabled flag. Set value to '1' to enable
# scheduled Telegram posting; default ships disabled so the cron is a
# no-op until the user opts in.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * www-data cd /var/www/html && /usr/local/bin/php cron/telegram_sender.php >> /var/www/html/logs/cron_stdout.log 2>&1
